# Discord Bot Files/discord-news.service
[Unit]
Description=Discord News Bot
After=network.target

[Service]
User=ubuntu
WorkingDirectory=/home/ubuntu
ExecStart=/usr/bin/python3 /home/ubuntu/discord-news-bot.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
